# launchd/com.todo-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.todo-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/einargudjonsson/personal/todo-system/.venv/bin/todo-sync</string>
        <string>all</string>
    </array>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>StandardOutPath</key>
    <string>/tmp/todo-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/todo-sync.err</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
